# profile.d/_go2dir

#!/bin/bash
local OPTIND OPTARG o usage debug direction dir

# map="
#    plugin        ~/AppData/Local/Google/Chrome/User Data/Default/Extensions   
# "

usage="
usage:

    $prog direction

directions: $map
"

debug=N
while getopts d o; do
    case $o in
        d) debug=Y ;;
        *) echo "unknown switch='$o'. $usage" >&2; return ;;
    esac
done

shift $((OPTIND - 1))

if [ $# -ne 1 ]; then
    echo "wrong number of arguments. $usage" >&2
    return
fi

direction="$1"

[  $debug = Y ] && set -x

dir=$(echo "$map" | grep -i " $direction " | perlfields 2 NF | sed 's/^[[:space:]]*//;s/[[:space:]]*$//')
if [ "X$dir" = "X" ]; then
    echo "unknown direction='$direction'. $usage" >&2
else
   # we run tptime to set time variables if dir has yyyy mm dd YYMD HH MM SS, case insensitive
   if echo "$dir" | grep -iqE '\b([pt]?yyyy|mm|dd|YYMD|HH|MM|SS)\b'; then
      # run tptime to set time variables
      tptime set
      [ $debug = Y ] && echo "tptime set"
   fi
   [ $debug = Y ] && echo "before eval dir='$dir'"

   # expand ~ and variables
   dir=$(eval echo "$dir") 
   [ $debug = Y ] && echo "after eval dir='$dir'"

   # double quotes to preserve spaces       
   cd "$dir"                
fi

[ $debug = Y ] && set +x
