# profile.d/mydir
local OPTIND OPTARG o quiet usage prog old new value var exact

prog=${FUNCNAME[1]}

usage="
usage:
   $prog pattern

   -e     exact name. default is a pattern.
   -l     list

example:
   $prog angular
   $prog -e lca

"

exact=N

while getopts el o; do
   case "$o" in
   e) exact=Y ;;
   l) ls -d "$MYBASE"/github/*
      return
      ;;
   *) echo "unknown switch. $usage" >&2
      return
      ;;
   esac
done

shift $((OPTIND - 1))

if [ $# -ne 1 ]; then
   echo "wrong number of args: expected 1, actual $#"
   echo "$usage"
   return 
fi

if [ $exact = Y ]; then
   cd "$MYBASE"/github/$1
else
   count=`/bin/ls -1 -d "$MYBASE"/github/$1*|wc -l`
   if [ $count -gt 1 ]; then
      /bin/ls -1l -d "$MYBASE"/github/$1*
      echo "ERROR: $prog: more than 1 matched. use '-e' to select one"
      return
   fi

   cd "$MYBASE"/github/$1*
fi
