# profile.d/rustenv
# for vs code
local found CARGO_DIR WABT_DIR usage prog

prog=${FUNCNAME[1]}

usage="

usage:

   $prog check
   $prog set

   set and check rust and wasm env

"

if [ $# -ne 1 ]; then
   echo "$prog: wrong number of args"
   echo "$usage"
   return
fi

action=$1

if [[ $UNAME =~ CYGWIN ]]; then
   CARGO_DIR=/cygdrive/c/users/$USERNAME/.cargo
   WABT_DIR=$(ls -1dt /cygdrive/c/users/$USERNAME/local/wabt* | tail -1)
elif [[ $UNAME =~ MINGW64_NT ]]; then
   CARGO_DIR=~/c/users/$USERNAME/.cargo
   WABT_DIR=$(ls -1dt /c/users/$USERNAME/local/wabt* | tail -1)
elif [[ $UNAME =~ Linux ]]; then
   CARGO_DIR=~/.cargo
   WABT_DIR=$(ls -1dt ~/local/wabt* | tail -1)
elif [[ $UNAME =~ Darwin ]]; then
   CARGO_DIR=~/.cargo
   WABT_DIR=$(ls -1dt ~/local/wabt* | tail -1)
else
   echo "unsupported UNAME='$UNAME'"
   return
fi

if [ $action = check ]; then
   echo ""
   echo "current installations: "
   (
      set -x
      /bin/ls -l "$CARGO_DIR"
      /bin/ls -l "$WABT_DIR"
   )
elif [ $action = set ]; then
   export PATH="$CARGO_DIR/bin:$WABT_DIR/bin:$PATH"
   reduce PATH
   (
      set -x
      which cargo
      which wasm2wat
   )
else
   echo "unknown action='$action'"
   usage
fi
